# docker/pi-hole/crontab
# Update Graviton's DB once a week.
0 0 * * 0 docker ps -f 'name=pihole' -f 'status=running' -f 'health=healthy' -q | xargs -I{} docker exec {} pihole -g
